# watershed.service
[Unit]
Description=Watershed catchment calculator
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/watershed
ExecStart=%h/watershed/venv/bin/python %h/watershed/sheed.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
